# vendor-camera.service
# =============================================================================
# 駄菓子屋さんカメラ - systemd ユーザーサービス
# -----------------------------------------------------------------------------
# ユーザーがログインしたときに自動で start_camera.sh を起動します。
# （ユーザーサービスはログイン時に開始され、ログアウト時に停止します）
#
# 設置・有効化の手順は ubuntu_autostart_setup.txt を参照してください。
# =============================================================================
[Unit]
Description=Dagashiya Camera (Flask + Vite) auto start
After=graphical-session.target network-online.target
Wants=network-online.target

[Service]
Type=simple

# ---- 環境に合わせて書き換えてください ----
# プロジェクトを置いた絶対パス
WorkingDirectory=%h/vendor_camera_ts
# 起動スクリプト（CAMERA_ID は必要に応じて変更）
Environment=CAMERA_ID=2
ExecStart=/usr/bin/env bash %h/vendor_camera_ts/start_camera.sh
# -----------------------------------------

# 異常終了したら自動で再起動
Restart=on-failure
RestartSec=5

# ログを journald に出力（journalctl --user -u vendor-camera で確認可能）
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
