# vendor-camera-dht.timer
# =============================================================================
# DHT 温湿度ロガー - systemd タイマー
# -----------------------------------------------------------------------------
# PC の時計を使って、毎時 0分,10分,20分,30分,40分,50分 に
# vendor-camera-dht.service（ワンショット測定）を起動します。
#
# 設置手順は ubuntu_autostart_setup.txt の「温湿度ロガー」章を参照。
# =============================================================================
[Unit]
Description=Run DHT logger every 10 minutes (:00 :10 :20 ...)

[Timer]
# 毎時 0,10,20,30,40,50 分の 00 秒に起動
OnCalendar=*-*-* *:00/10:00
# 起動を逃した分（スリープ等）は復帰後に 1 回だけ実行
Persistent=true
AccuracySec=1s

[Install]
WantedBy=timers.target
