# vendor-camera-dht.service
# =============================================================================
# DHT 温湿度ロガー - systemd サービス（ワンショット）
# -----------------------------------------------------------------------------
# vendor-camera-dht.timer から起動され、1 回測定して SQLite に記録します。
# シリアル(/dev/ttyACM*)アクセスのため、ユーザーを dialout グループに入れるか
# root で実行してください（下記は root 実行・system サービス例）。
#
# 設置手順は ubuntu_autostart_setup.txt の「温湿度ロガー」章を参照。
# =============================================================================
[Unit]
Description=DHT temperature/humidity logger (one shot)
After=multi-user.target

[Service]
Type=oneshot

# ---- 環境に合わせて絶対パスを書き換えてください ----
ExecStart=/usr/bin/python3 /home/USER/vendor_camera_ts/dht_logger.py
# シリアルポートや DB の場所を変える場合は Environment で上書き可能
# Environment=DHT_SERIAL_PORT=/dev/ttyACM0
# Environment=DHT_DB_PATH=/home/USER/vendor_camera_ts/output/sensor.db
# ----------------------------------------------------

StandardOutput=journal
StandardError=journal
